# systemd/hubitat-event.service
[Unit]
Description=Hubitat Event Polling Daemon (Docker)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/hubitat-event
ExecStartPre=-/usr/bin/docker rm -f hubitat-event
ExecStartPre=/bin/bash -c 'if [ "${AUTO_PULL}" = "1" ]; then /usr/bin/docker pull "${IMAGE}"; fi'
ExecStart=/usr/bin/docker run --rm \
  --name hubitat-event \
  $DOCKER_EXTRA_ARGS \
  -e EVENT_CFG_JSON=${EVENT_CFG_JSON} -e POLL_INTERVAL=${POLL_INTERVAL} \
  -v /opt/hubitat/cfg:/home/admin/cfg \
  -v /opt/hubitat/cache:/home/admin/cache \
  -v /home/admin/.ssh:/home/admin/.ssh:ro \
  ${IMAGE}
ExecStop=/usr/bin/docker stop hubitat-event
Restart=always
RestartSec=5s
TimeoutStopSec=30s
KillMode=control-group
OOMScoreAdjust=-400
LimitNOFILE=65536

[Install]
WantedBy=hubitat.target
